# bwv_sync/bwv-sync.cron
35 5 * * * root /usr/local/bin/looplijsten_bwv_sync.sh
26 16 * * * root /usr/local/bin/looplijsten_bwv_sync.sh
